# config/task.cron
*/1 * * * * cd /var/www/html/ && nice -n -20 /usr/local/bin/php adm command job --yes > /var/www/html/storage/logs/job.log 2>&1
